# ultimate-sovereign-rpc-master-bundle/sovereign-rpc-multi-os-support/sovereign-rpc-centos/systemd/caddy.service
[Unit]
Description=Caddy Reverse Proxy
After=network.target docker.service
Requires=docker.service

[Service]
ExecStart=/usr/bin/docker compose -f /opt/sovereign-rpc/docker-compose.yml up caddy
Restart=always
User=root

[Install]
WantedBy=multi-user.target