# ultimate-sovereign-rpc-master-bundle/sovereign-rpc-multi-os-support/sovereign-rpc-ubuntu/systemd/prometheus.service
[Unit]
Description=Prometheus Monitoring
After=network.target docker.service
Requires=docker.service

[Service]
ExecStart=/usr/bin/docker compose -f /opt/sovereign-rpc/docker-compose.yml up prometheus
Restart=always
User=root

[Install]
WantedBy=multi-user.target